# scripts/com.dooley.mailtriage.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.dooley.mailtriage.daily</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/dooley/Documents/GithubClone/MailTriage/scripts/run_daily_mailtriage.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>5</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/mailtriage-daily.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mailtriage-daily.err.log</string>
  </dict>
</plist>
